# build/deb.d/telegraf-docker-sd.service
[Unit]
Description=Telegraf Docker service discovery daemon
After=network.target

[Service]
Type=simple
User=telegraf
WorkingDirectory=/etc/telegraf
ExecStart=/usr/local/bin/telegraf-docker-sd
EnvironmentFile=/etc/telegraf/telegraf-docker-sd.env
Restart=on-abort

[Install]
WantedBy=multi-user.target
