# docker/etc/services.d/telegraf-docker-sd/run
#!/usr/bin/with-contenv sh

s6-envuidgid telegraf /usr/local/bin/telegraf-docker-sd
